# mysql/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/run/supervisord.pid
#loglevel=debug

[program:mysql]
command=/etc/init.d/mysqld start
autorestart=true
priority=1
